# package/qca-nss-drv/files/qca-nss-drv.init
#!/bin/sh /etc/rc.common
#
# Copyright (c) 2014, 2016 The Linux Foundation. All rights reserved.

START=70

start() {

	irq_nss_rps=`grep -m3 nss /proc/interrupts | cut -d ':' -f 1 | tail -n1 | tr -d ' '`
	[ -n "$irq_nss_rps" ] && echo 2 > /proc/irq/$irq_nss_rps/smp_affinity

	# Enable NSS RPS
	sysctl -w dev.nss.n2hcfg.rps=1 >/dev/null 2>/dev/null

}
