# target/linux/ar71xx/CUS227/base-files/etc/init.d/skiftahostname
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=11

start() {

	local hostname="$(uci_get system.@system[0].hostname)"
	if [ "$hostname" == "Unconfigured-Skifta" ]; then
		local macaddr="$(cat /sys/class/ieee80211/phy0/macaddress | tr 'A-Z' 'a-z')"
		local ssidString="Skifta_${macaddr:6:12}"
		echo "$ssidString" > /proc/sys/kernel/hostname
	else
		echo "$hostname" > /proc/sys/kernel/hostname
	fi
}

stop() {
	true
}
